# provisioning/all-in-one/files/supervisor.golang.conf
[program:golang]
directory=/home/isucon/webapp/golang
command=/home/isucon/env.sh ./app
user=isucon
stdout_logfile=/tmp/isucon.golang.log
stderr_logfile=/tmp/isucon.golang.log
autostart=false
