# provisioning/image/files/supervisor.node.conf
[program:node]
directory=/home/isucon/webapp/node
command=/home/isucon/env.sh node app.js
user=isucon
stdout_logfile=/tmp/isucon.node.log
stderr_logfile=/tmp/isucon.node.log
autostart=false
